# tinllm.service
[Unit]
Description=TinLLM MoE Inference Engine Service
After=network.target

[Service]
Type=simple
User=user
WorkingDirectory=/home/user/Tinllm
ExecStart=/usr/bin/python3 /home/user/Tinllm/openai_server.py
Restart=always
RestartSec=3
Environment=OMP_NUM_THREADS=8
Environment=DIRECT=1
Environment=URING=1

[Install]
WantedBy=multi-user.target
